# bin/elsystemd.service
[Unit]
Description=Enter Leave System Daemon
After=network.target local-fs.target
Requires=network.target local-fs.target

[Service]
User=geek
Group=geek
EnvironmentFile=/home/pi/EnterLeaveSystem/env
WorkingDirectory=/home/pi/EnterLeaveSystem
ExecStart=/home/pi/EnterLeaveSystem/elsystem
ExecStop=/bin/kill -SIGINT $MAINPID
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target
